# de/.local/sv/wl-gammarelay-rs/run
#!/bin/sh
# Gamma/brightness D-Bus service (rs.wl-gammarelay). Needs the compositor's
# wayland socket + the session bus. runit respawns it after a dwl restart (it
# exits when the old display goes away) and across the boot monitor-bringup
# race, so it's always back up — redshift.sh / brightness.sh just drive it.
bus="${DBUS_SESSION_BUS_ADDRESS#unix:path=}"
wl="${XDG_RUNTIME_DIR}/${WAYLAND_DISPLAY}"
while [ ! -S "$bus" ] || [ ! -e "$wl" ]; do sleep 0.2; done
exec wl-gammarelay-rs
