# de/.local/sv/syncthing/run
#!/bin/sh
# --no-restart: let runsv supervise; don't let syncthing fork its own monitor
exec syncthing serve --no-browser --no-restart
