# de/.local/sv/dbus/run
#!/bin/sh
# Persistent user D-Bus session bus. autostart.sh exports DBUS_SESSION_BUS_ADDRESS
# before starting runsvdir, so this service and all siblings inherit the same address.
xrd="${XDG_RUNTIME_DIR:-/run/user/$(id -u)}"
rm -f "$xrd/bus"
exec dbus-daemon --session --address="unix:path=$xrd/bus" --nofork --nopidfile
