# de/.local/sv/wireplumber/run
#!/bin/sh
xrd="${XDG_RUNTIME_DIR:-/run/user/$(id -u)}"
while [ ! -e "$xrd/pipewire-0" ] || [ ! -S "$xrd/bus" ]; do sleep 0.2; done
exec wireplumber
