# de/.local/sv/cliphist-text/run
#!/bin/sh
# Clipboard-history watcher for text (wl-paste -> cliphist store). Needs a live
# wayland socket, which may be stale if runsvdir outlives a compositor restart
# — wait for it instead of trusting the inherited WAYLAND_DISPLAY.
wl=$(ls "${XDG_RUNTIME_DIR}"/wayland-[0-9] 2>/dev/null | head -1)
while [ -z "$wl" ]; do sleep 0.2; wl=$(ls "${XDG_RUNTIME_DIR}"/wayland-[0-9] 2>/dev/null | head -1); done
export WAYLAND_DISPLAY="${wl##*/}"
exec wl-paste --type text --watch cliphist store
