# de/.local/sv/mako/run
#!/bin/sh
# Notification daemon (mako). Needs the session bus + a live wayland socket;
# both may be absent/stale if runsvdir outlives a compositor restart, so wait
# for them rather than trusting the inherited env.
bus="${DBUS_SESSION_BUS_ADDRESS#unix:path=}"
while [ ! -S "$bus" ]; do sleep 0.2; done
wl=$(ls "${XDG_RUNTIME_DIR}"/wayland-[0-9] 2>/dev/null | head -1)
while [ -z "$wl" ]; do sleep 0.2; wl=$(ls "${XDG_RUNTIME_DIR}"/wayland-[0-9] 2>/dev/null | head -1); done
export WAYLAND_DISPLAY="${wl##*/}"
exec mako
